# docker/root/etc/laravel-echo-server/laravel-echo-server.ini
[program:laravel-echo-server]
command=/usr/bin/laravel-echo-server start --force
stdout_logfile = /var/log/laravel-echo-server
stdout_events_enabled = true
